# roles/touchscreen/files/turn_off_display.service
[Unit]
Description=Turn off display

[Service]
Type=oneshot
ExecStart=/usr/local/bin/turn_off_display.sh
