# roles/touchscreen/files/turn_on_display.timer
[Unit]
Description=Turn on display at scheduled times

[Timer]
OnCalendar=*-*-* 05:30:00,16:30:00
Persistent=true

[Install]
WantedBy=timers.target
