# roles/touchscreen/files/turn_on_display.service
[Unit]
Description=Turn on display

[Service]
Type=oneshot
ExecStart=/usr/local/bin/turn_on_display.sh
