# roles/touchscreen/templates/turn_off_display_after_touch.service.j2
[Unit]
Description=Turn off HDMI display after touchscreen wake

[Service]
ExecStart=/usr/local/bin/turn_off_display.sh
