# roles/touchscreen/templates/turn_on_display_pm.service.j2
[Unit]
Description=Turn on HDMI display (PM)

[Service]
ExecStart=/usr/local/bin/turn_on_display.sh
