# roles/touchscreen/templates/turn_off_display.service.j2
[Unit]
Description=Turn off HDMI display

[Service]
ExecStart=/usr/local/bin/turn_off_display.sh
