# roles/touchscreen/templates/turn_on_display_touch.service.j2
[Unit]
Description=Turn on HDMI display (Touch)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/turn_on_display.sh
ExecStartPost=/bin/systemctl start turn_off_display_after_touch.timer
