# roles/touchscreen/templates/turn_off_display_am.timer.j2
[Unit]
Description=Turn off HDMI display at 09:00

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
